# content/etc/s6-overlay/s6-rc.d/svc-postfix/run
#!/command/with-contenv bash

echo "Starting postfix"
/usr/sbin/postfix start-fg